# src/argo-voms-htpasswd.cron

# Run argo-voms-htpasswd if enabled via the argo-voms-htpasswd init.d script.
0 */6 * * * root [ -f /var/lock/subsys/argo-voms-htpasswd ] && sleep `expr $RANDOM \% 100` && ( /usr/sbin/argo-voms-htpasswd 2>&1 | gawk '{print strftime(), ":", $0}'  >> /var/log/argo-voms-htpasswd.log  )



